# supervisord/statsite.conf
[program:statsite]
command = statsite -f /etc/statsite.ini
user = root
autostart = true
startsecs = 5
autorestart = true
user = root
redirect_stderr = true
